# deb/DEBIAN/postinst
cd /home/plica/mgmt && pip3 install -r requeriments.txt
chmod 777 /opt/sensoruba/data/*
chmod 777 /tmp/sensor-watchdog
systemctl daemon-reload
systemctl enable ubastart.service
systemctl start ubastart.service
